# m-observe-overseer/install/m-observe-host@.service
[Unit]
Description=M-OBSERVE Host (Backend + WebUI)
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/m-observe/backend
ExecStart=/opt/m-observe/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 3501
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
